# tesla_http_proxy/rootfs/etc/s6-overlay/s6-rc.d/webui/run
#!/command/with-contenv bashio

# read options
# you can pass in these variables if running without supervisor
if [ -n "${HASSIO_TOKEN:-}" ]; then
  DOMAIN="$(bashio::config 'domain')"; export DOMAIN
  DEBUG="$(bashio::config 'debug')"; export DEBUG
fi

python3 /app/webui.py
